# scripts/com.myloggy.synclog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.myloggy.synclog</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/takehiroichikawa/myloggy/scripts/sync_to_notion.py</string>
    </array>
    <key>StartInterval</key>
    <integer>10800</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/takehiroichikawa/myloggy/logs/launchd-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/takehiroichikawa/myloggy/logs/launchd-stderr.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/takehiroichikawa/myloggy</string>
</dict>
</plist>
